# deploy/supervisor/abbev-bunny-worker.conf
; ---------------------------------------------------------------------------
; Worker Supervisor — transfert des vidéos vers Bunny (queue dédiée "bunny")
;
; Installe ce fichier dans /etc/supervisor/conf.d/ puis :
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start abbev-bunny-worker:*
;
; Adapte : le chemin du projet, l'utilisateur (www-data en général) et, si
; besoin, le binaire php. NE PAS mettre --timeout court : un PUT de plusieurs
; Go est légitime (le retry_after de la connexion "bunny" est à 7200s).
; ---------------------------------------------------------------------------
[program:abbev-bunny-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/abbev/artisan queue:work bunny --queue=bunny --timeout=0 --tries=3 --sleep=3 --backoff=30
directory=/var/www/abbev
autostart=true
autorestart=true
stopwaitsecs=7200
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/abbev/storage/logs/bunny-worker.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stopsignal=TERM
